# contrib/ugpibd.service
[Unit]
Description=Agilent/Keysight 82357B GPIB daemon (Prologix-compatible)
After=network.target

[Service]
ExecStart=/usr/local/bin/ugpibd
Restart=on-failure
RestartSec=3
User=ugpibd
Group=plugdev
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
